# recipes-core/systemd/systemd-udev/50-ppp0.rules
ACTION=="add", SUBSYSTEM=="net", ENV{INTERFACE}=="ppp0", RUN+="/bin/systemctl --no-block start systemd-time-wait-sync.service"