# recipes-core/systemd/systemd-udev/50-ppp0.rules
ACTION=="remove", SUBSYSTEM=="net", ENV{INTERFACE}=="ppp0", RUN+="/bin/systemctl --no-block stop systemd-networkd-wait-online@ppp0.service"