# recipes-core/systemd/systemd-udev/50-sda.rules
ACTION=="add", SUBSYSTEM=="block", ENV{DEVNAME}=="/dev/sda", RUN+="/bin/systemctl --no-block start fwu-usb.service"