# recipes-core/systemd/systemd-units/ble.target
[Unit]
Description=ble
After=check.target
Wants=check.target
