# recipes-core/systemd/systemd-units/check.target
[Unit]
Description=check
After=drive.target
Wants=drive.target
