# recipes-core/systemd/systemd-units/chk-usb.service
[Unit]
Description=chk-usb

[Service]
Type=simple
ExecStart=/etc/scripts/chk-usb.sh
