# recipes-core/systemd/systemd-units/communication.target
[Unit]
Description=communication
After=check.target
Wants=check.target
