# recipes-core/systemd/systemd-units/communication.target
[Unit]
Description=Communication
After=drive.target
Wants=drive.target
