# recipes-core/systemd/systemd-units/communication.target
[Unit]
Description=communication
After=drive.target
Wants=drive.target
