# recipes-core/systemd/systemd-units/debug.target
[Unit]
Description=debug
After=update.target
Wants=update.target
