# recipes-core/systemd/systemd-units/debug.target
[Unit]
Description=debug
After=late-init.target
Wants=late-init.target
