# recipes-core/systemd/systemd-units/debug.target
[Unit]
Description=debug
After=location.target
Wants=location.target
