# recipes-core/systemd/systemd-units/debug.target
[Unit]
Description=debug
After=communication.target
Wants=communication.target
