# recipes-core/systemd/systemd-units/drive.target
[Unit]
Description=drive
DefaultDependencies=no
Wants=early-init.target
After=early-init.target
