# recipes-core/systemd/systemd-units/drive.target
[Unit]
Description=drive
After=graphical.target
Wants=graphical.target
