# recipes-core/systemd/systemd-units/drive.target
[Unit]
Description=Drive
After=graphical.target
Wants=graphical.target
