# recipes-core/systemd/systemd-units/drive.target
[Unit]
Description=drive
DefaultDependencies=no
After=sysinit.target
Wants=sysinit.target
