# recipes-core/systemd/systemd-units/early-init.target
[Unit]
Description=Early init target for everything which needs to run for medusa
Wants=drive.target
