# recipes-core/systemd/systemd-units/fwu-usb-chk.service
[Unit]
Description=fwu-usb-chk

[Service]
Type=simple
ExecStart=/etc/scripts/fwu-usb-chk.sh
