# recipes-core/systemd/systemd-units/fwu-usb.service
[Unit]
Description=fwu-usb
After=mnt-rfs.service led.service drive.target
Requires=mnt-rfs.service
Wants=led.service drive.target

[Service]
Type=simple
ExecStart=/etc/scripts/fwu-usb.sh
KillSignal=SIGHUP
