# recipes-core/systemd/systemd-units/fwu.service
[Unit]
Description=fwu

[Service]
Type=oneshot
ExecStart=/etc/scripts/fwu.sh
RemainAfterExit=yes

[Install]
WantedBy=basic.target
