# recipes-core/systemd/systemd-units/gps.target
[Unit]
Description=gps
After=ble.target
Wants=ble.target
