# recipes-core/systemd/systemd-units/gsm.target
[Unit]
Description=gsm
After=gps.target
Wants=gps.target
