# recipes-core/systemd/systemd-units/hostname-set.service
[Unit]
Description=hostname-set
RequiresMountsFor=/mnt/data
DefaultDependencies=no
Before=rsyslog.service

[Service]
Type=oneshot
ExecStart=/bin/dash -c "/bin/hostname -F /mnt/data/hostname || true"
RemainAfterExit=yes

[Install]
WantedBy=drive.target
