# recipes-core/systemd/systemd-units/late-init.target
[Unit]
Description=Late init target for everything which needs to run eventually
After=gsm.target
Wants=gsm.target
