# recipes-core/systemd/systemd-units/led.service
[Unit]
Description=led
DefaultDependencies=no
After=can0.service

[Service]
Type=oneshot
ExecStart=/etc/scripts/led.sh start
ExecStop=/etc/scripts/led.sh stop
RemainAfterExit=yes

[Install]
WantedBy=drive.target
