# recipes-core/systemd/systemd-units/location.target
[Unit]
Description=location
After=communication.target
Wants=communication.target
