# recipes-core/systemd/systemd-units/log-usb.service
[Unit]
Description=log-usb
After=mnt-log.service
Wants=mnt-log.service led.service

[Service]
Type=oneshot
ExecStart=/etc/scripts/log-usb.sh
