# recipes-core/systemd/systemd-units/mnt-data.mount
[Unit]
Description=mnt-data

[Mount]
What=/dev/ubi0_2
Where=/mnt/data
Type=ubifs
