# recipes-core/systemd/systemd-units/mnt-rfs.service
[Unit]
Description=mnt-rfs
After=drive.target

[Service]
Type=oneshot
ExecStart=/etc/scripts/mnt-rfs.sh start
ExecStop=/etc/scripts/mnt-rfs.sh stop
RemainAfterExit=yes
