# recipes-core/systemd/systemd-units/mnt-rfs.service
[Unit]
Description=Mount the inactive rfs partition

[Service]
Type=oneshot
ExecStart=/etc/scripts/mnt-rfs.sh start
ExecStop=/etc/scripts/mnt-rfs.sh stop
RemainAfterExit=yes

[Install]
WantedBy=basic.target
