# recipes-core/systemd/systemd-units/started.target
[Unit]
Description=started
After=debug.target multi-user.target
Wants=debug.target multi-user.target
