# recipes-core/systemd/systemd-units/update.target
[Unit]
Description=update
After=communication.target
Wants=communication.target
