# recipes-core/systemd/systemd-units/watchdog-ppp0.service
[Unit]
Description=watchdog-ppp0

[Service]
Type=simple
ExecStart=/etc/scripts/watchdog-ppp0.sh